# base/root/etc/s6-overlay/s6-rc.d/init-prepare-directories/run
#!/usr/bin/with-contenv bash

echo "Hello world!"