# nodejs/root/etc/s6-overlay/s6-rc.d/node-install-packages/up
/etc/s6-overlay/s6-rc.d/node-install-packages/run