# nodejs/root/etc/s6-overlay/s6-rc.d/node-run-build-script/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid dolphin /etc/s6-overlay/s6-rc.d/node-run-build-script/rootless-run
