# nodejs/root/etc/s6-overlay/s6-rc.d/node-run-build-script/up
/etc/s6-overlay/s6-rc.d/node-run-build-script/run