# python/root/etc/s6-overlay/s6-rc.d/python-init-done/run
#!/usr/bin/with-contenv bash

echo "[+] Python init done."
echo