# base/root/etc/s6-overlay/s6-rc.d/init-done/run
#!/usr/bin/with-contenv bash

echo "[+] Environment setup completed"
echo