# nodejs/root/etc/s6-overlay/s6-rc.d/node-init-done/run
#!/usr/bin/with-contenv bash

echo "[+] Node.js init done."
echo