# nodejs/root/etc/s6-overlay/s6-rc.d/node-init-done/up
/etc/s6-overlay/s6-rc.d/node-init-done/run